# arb.timer
[Unit]
Description=Automatic Robust Backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=graphical-session.target

