# arb.service
[Unit]
Description=Automatic Robust Backup
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/arb
