# arbie.service
[Unit]
Description=Automatic Robust Backup
Wants=default.target

[Service]
Type=oneshot
ExecStart=arbie
